# hyprland/profiles/arch-plasma/runtime/autostart/monochrome-media-visualizer.desktop
[Desktop Entry]
Type=Application
Name=KDE Rice Preview Placeholder 2
Exec=true
Terminal=false
X-KDE-AutostartScript=true
